# roles/deploy_sample_webapp/files/web_app_supervisord.conf
[supervisord]
[program:gunicorn]
directory=/srv/web_app/
command=gunicorn app:app -b 0.0.0.0:8000 -w 5